# backend/config/supervisord/airplay-multi-instance.ini
# Multi-Instance AirPlay Configuration
# 3 independent shairport-sync instances with separate lifecycle managers
# Each endpoint can be used simultaneously by different devices/rooms

# ============================================================================
# AIRPLAY INSTANCE 1 (Primary - port 5000, D-Bus enabled for control)
# ============================================================================

[program:shairport-sync-1]
command=/usr/local/bin/shairport-sync --configfile=/app/config/shairport-sync-1.conf
user=root
directory=/app
environment=HOME="/app"
priority=40
autostart=true
autorestart=true
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisord/shairport-1.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/supervisord/shairport-1_err.log
stderr_logfile_maxbytes=10MB

[program:airplay-1-fifo-keeper]
command=/bin/bash /app/scripts/airplay-multi-fifo-keeper.sh 1
user=snapcast
directory=/app
priority=20
autostart=true
autorestart=true
startsecs=5
startretries=999
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisord/airplay-1-fifo-keeper.log
stderr_logfile=/var/log/supervisord/airplay-1-fifo-keeper_err.log

[program:airplay-1-lifecycle-manager]
command=/usr/bin/python3 /app/scripts/stream-lifecycle-manager.py --instance-id 1
user=snapcast
directory=/app
priority=35
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisord/airplay-1-lifecycle.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisord/airplay-1-lifecycle_err.log
stderr_logfile_maxbytes=10MB
environment=PYTHONUNBUFFERED="1"

# ============================================================================
# AIRPLAY INSTANCE 2 (Secondary - port 7000, D-Bus disabled)
# ============================================================================

[program:shairport-sync-2]
command=/usr/local/bin/shairport-sync --configfile=/app/config/shairport-sync-2.conf
user=root
directory=/app
environment=HOME="/app"
priority=40
autostart=true
autorestart=true
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisord/shairport-2.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/supervisord/shairport-2_err.log
stderr_logfile_maxbytes=10MB

[program:airplay-2-fifo-keeper]
command=/bin/bash /app/scripts/airplay-multi-fifo-keeper.sh 2
user=snapcast
directory=/app
priority=20
autostart=true
autorestart=true
startsecs=5
startretries=999
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisord/airplay-2-fifo-keeper.log
stderr_logfile=/var/log/supervisord/airplay-2-fifo-keeper_err.log

[program:airplay-2-lifecycle-manager]
command=/usr/bin/python3 /app/scripts/stream-lifecycle-manager.py --instance-id 2
user=snapcast
directory=/app
priority=35
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisord/airplay-2-lifecycle.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisord/airplay-2-lifecycle_err.log
stderr_logfile_maxbytes=10MB
environment=PYTHONUNBUFFERED="1"

# ============================================================================
# AIRPLAY INSTANCE 3 (Tertiary - port 8000, D-Bus disabled)
# ============================================================================

[program:shairport-sync-3]
command=/usr/local/bin/shairport-sync --configfile=/app/config/shairport-sync-3.conf
user=root
directory=/app
environment=HOME="/app"
priority=40
autostart=true
autorestart=true
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisord/shairport-3.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/supervisord/shairport-3_err.log
stderr_logfile_maxbytes=10MB

[program:airplay-3-fifo-keeper]
command=/bin/bash /app/scripts/airplay-multi-fifo-keeper.sh 3
user=snapcast
directory=/app
priority=20
autostart=true
autorestart=true
startsecs=5
startretries=999
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisord/airplay-3-fifo-keeper.log
stderr_logfile=/var/log/supervisord/airplay-3-fifo-keeper_err.log

[program:airplay-3-lifecycle-manager]
command=/usr/bin/python3 /app/scripts/stream-lifecycle-manager.py --instance-id 3
user=snapcast
directory=/app
priority=35
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisord/airplay-3-lifecycle.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisord/airplay-3-lifecycle_err.log
stderr_logfile_maxbytes=10MB
environment=PYTHONUNBUFFERED="1"
