# backend/config/supervisord/auto-switch.ini
# Auto-switch service: monitors stream activity and switches groups/snapclient target
[program:auto-switch]
command=python3 /app/scripts/auto-switch-service.py
user=root
priority=35
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisord/auto-switch.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisord/auto-switch_err.log
stderr_logfile_maxbytes=10MB
