# backend/config/supervisord/bluetooth-fifo-keeper.ini
# Bluetooth FIFO Keeper
# Keeps Bluetooth FIFO pipe open when no Snapcast stream exists
# Prevents bluealsa-aplay from blocking when writing to FIFO
[program:bluetooth-fifo-keeper]
command=/bin/bash /app/scripts/bluetooth-fifo-keeper.sh /tmp/bluetooth-fifo
user=snapcast
directory=/app
priority=26
autostart=false
autorestart=true
startsecs=3
startretries=999
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisord/bluetooth-fifo-keeper.log
stderr_logfile=/var/log/supervisord/bluetooth-fifo-keeper_err.log
