# backend/config/supervisord/dlna-fifo-keeper.ini
# DLNA FIFO Keeper
# Prevents gmrender-resurrect from blocking on FIFO writes when no Snapcast stream exists
# Reads and discards data from FIFO when stream is not active
[program:dlna-fifo-keeper]
command=/bin/bash /app/scripts/dlna-fifo-keeper.sh
user=snapcast
directory=/app
priority=45
autostart=false
autorestart=true
startsecs=3
startretries=999
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisord/dlna-fifo-keeper.log
stderr_logfile=/var/log/supervisord/dlna-fifo-keeper_err.log
