# backend/config/supervisord/federation.ini
[program:federation]
command=python3 -m federation.service
directory=/app/scripts
user=snapcast
priority=50
autostart=true
autorestart=true
startsecs=1
startretries=3
stderr_logfile=/var/log/supervisord/federation-stderr.log
stdout_logfile=/var/log/supervisord/federation-stdout.log
environment=PYTHONUNBUFFERED=1
