# backend/config/supervisord/fifo-keeper.ini
# FIFO Keeper - Keeps AirPlay FIFO open for reading
# This prevents shairport-sync from blocking when no Snapcast stream exists
[program:airplay-fifo-keeper]
command=/bin/sh /app/scripts/fifo-keeper.sh /tmp/snapfifo
user=snapcast
directory=/app
priority=20
autostart=true
autorestart=true
startsecs=2
startretries=999
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisord/fifo-keeper.log
stdout_logfile_maxbytes=1MB
stderr_logfile=/var/log/supervisord/fifo-keeper_err.log
stderr_logfile_maxbytes=1MB
