# backend/config/supervisord/init-services.ini
# Service Initialization Script
# Runs once at startup to start services based on settings.json
# Priority 10 ensures this runs before other services that depend on settings
[program:init-services]
command=/usr/bin/python3 /app/scripts/init-services-from-settings.py
user=snapcast
directory=/app
priority=10
autostart=true
autorestart=false
startsecs=0
startretries=0
exitcodes=0
stdout_logfile=/var/log/supervisord/init-services.log
stderr_logfile=/var/log/supervisord/init-services_err.log
environment=PYTHONUNBUFFERED="1"
