# backend/config/supervisord/metadata-debug-server.ini
[program:metadata-debug-server]
command=/usr/bin/python3 -u /app/scripts/metadata-debug-server.py
user=snapcast
directory=/app
priority=60
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisord/metadata-debug-server.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisord/metadata-debug-server_err.log
stderr_logfile_maxbytes=10MB
environment=PYTHONUNBUFFERED="1"
