# backend/config/supervisord/mosquitto.ini
[program:mosquitto]
command=/usr/sbin/mosquitto -c /app/config/mosquitto.conf
user=mosquitto
directory=/app
priority=5
autostart=true
autorestart=true
startsecs=3
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisord/mosquitto.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisord/mosquitto_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
