# backend/config/supervisord/plexamp-stream-lifecycle-manager.ini
[program:plexamp-stream-lifecycle-manager]
command=/usr/bin/python3 /app/scripts/plexamp-stream-lifecycle-manager.py
user=snapcast
directory=/app
priority=25
autostart=true
autorestart=true
startsecs=5
startretries=999
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisord/plexamp-stream-lifecycle-manager.log
stderr_logfile=/var/log/supervisord/plexamp-stream-lifecycle-manager_err.log
environment=PYTHONUNBUFFERED="1"
