# backend/config/supervisord/services.ini
[program:dbus]
command=/bin/sh -c 'rm -f /var/run/dbus/pid && dbus-daemon --system --nofork --nopidfile'
autostart=true
autorestart=true
startsecs=3
startretries=5
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:avahi]
command=/bin/sh -c 'sleep 2 && avahi-daemon --no-drop-root'
autostart=true
autorestart=true
startsecs=3
startretries=5
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:shairport-sync]
command=/bin/sh -c 'if [ "%(ENV_AIRPLAY_CONFIG_ENABLED)s" = "1" ]; then sleep 3 && shairport-sync -c /app/config/shairport-sync.conf; else sleep infinity; fi'
autostart=true
autorestart=true
startsecs=5
startretries=5
priority=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:snapserver]
command=/bin/sh -c 'sleep 4 && snapserver -c /app/config/snapserver.conf --server.datadir=/app/data'
autostart=true
autorestart=true
startsecs=5
startretries=5
priority=40
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:snapclient]
command=/bin/sh -c 'if [ "%(ENV_SNAPCLIENT_ENABLED)s" = "1" ]; then sleep 8 && snapclient -h %(ENV_SNAPCLIENT_HOST)s -s %(ENV_SNAPCLIENT_SOUNDCARD)s --player alsa --latency %(ENV_SNAPCLIENT_LATENCY)s; else sleep infinity; fi'
autostart=true
autorestart=true
startsecs=10
startretries=5
priority=50
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true