# backend/config/supervisord/snapserver-avahi-restart.ini
# One-time service to restart snapserver after Avahi is ready
# This ensures zeroconf advertising works on initial container startup
[program:snapserver-avahi-restart]
command=/bin/bash /app/scripts/snapserver-avahi-restart.sh
user=root
directory=/app
priority=150
autostart=true
autorestart=false
startsecs=0
startretries=1
stdout_logfile=/var/log/supervisord/snapserver-avahi-restart.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisord/snapserver-avahi-restart_err.log
stderr_logfile_maxbytes=10MB
