# backend/config/supervisord/stream-manager.ini
# Stream Manager - Dynamic stream lifecycle management
# Monitors stream activity and manages idle timeouts
[program:stream-manager]
command=/usr/bin/python3 -u /app/scripts/stream-manager.py
user=snapcast
directory=/app
priority=35
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisord/stream-manager.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisord/stream-manager_err.log
stderr_logfile_maxbytes=10MB
environment=PYTHONUNBUFFERED="1",IDLE_TIMEOUT_MINUTES="%(ENV_IDLE_TIMEOUT_MINUTES)s"
