# backend/config/supervisord/supervisord.ini
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log


[program:dbus]
command=/bin/sh -c 'mkdir -p /var/run/dbus && dbus-daemon --system --nofork --nopidfile'
autostart=true
autorestart=true
startsecs=3
startretries=5
priority=10

[program:avahi]
command=/bin/sh -c 'sleep 1; avahi-daemon --no-drop-root'
autostart=true
autorestart=true
startsecs=3
startretries=5
priority=20

[program:shairport-sync]
command=/bin/sh -c 'sleep 2; shairport-sync -c /app/config/shairport-sync.conf'
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
startsecs=3
startretries=5
priority=25

[program:airplay-metadata]
command=/bin/bash /app/process-airplay-metadata.sh
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
startsecs=5
startretries=5
priority=26

[program:playback-controls]
command=/bin/bash /app/handle-playback-controls.sh
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
startsecs=5
startretries=5
priority=27

[program:snapserver]
command=/bin/sh -c 'sleep 3; snapserver -c /app/config/snapserver.conf --server.datadir=/app/data'
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
startsecs=3
startretries=5
priority=30
